# deploy/stackalpha-celery-worker.service
[Unit]
Description=StackAlpha Celery Worker
After=network.target redis.service postgresql.service
Wants=redis.service

[Service]
Type=forking
User=alpha
Group=alpha
WorkingDirectory=/home/alpha/stackalpha
Environment="PATH=/home/alpha/stackalpha/.venv/bin:/usr/local/bin:/usr/bin"
EnvironmentFile=/home/alpha/stackalpha/.env

ExecStart=/home/alpha/stackalpha/.venv/bin/celery \
    -A app.workers.celery_app worker \
    --loglevel=info \
    --concurrency=4 \
    --max-tasks-per-child=100 \
    --logfile=/var/log/stackalpha/celery-worker.log \
    --pidfile=/var/run/stackalpha/celery-worker.pid \
    --detach

ExecStop=/bin/kill -s TERM $MAINPID
ExecReload=/bin/kill -s HUP $MAINPID

PIDFile=/var/run/stackalpha/celery-worker.pid
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
